# device/testing/soc-unisoc-ums9230/99-sprd-watchdog.rules
ACTION=="add", SUBSYSTEM=="rpmsg", ATTR{name}=="pm_ctl0", RUN+="/bin/sh -c 'echo -n watchdog on > /dev/$name'"
