# device/testing/device-mecha-comet-m/bluetooth/mc-brcm-bt-enable.initd
#!/sbin/openrc-run
name="Mecha Comet brcm bluetooth enabler"
supervisor=supervise-daemon
command="/usr/bin/hciattach"
command_args="/dev/ttymxc0 bcm43xx 921600"
