# main/postmarketos-usb-moded/rootfs-etc-init.d-usb-moded-developer-mode.initd
#!/sbin/openrc-run

description="USB Moded Developer Mode"
command="usb-moded-developer-mode"
command_args="up"

stop() {
    usb-moded-developer-mode down
}
